# root/etc/services.d/haproxy/run
#!/usr/bin/with-contenv bash
CONFIG="/config/haproxy.cfg"
haproxy -f $CONFIG -c
haproxy -q -f $CONFIG -sf $(pidof haproxy)
