# services/haproxy/run
#!/bin/bash
CONFIG="/config/haproxy.cfg"

/usr/sbin/haproxy -f $CONFIG -c
/usr/sbin/haproxy -q -f $CONFIG -sf $(pidof haproxy)
