# defaults/letsencrypt.cron
37 13 */15 * *	root	exec /app/letsencrypt.sh
