# root/etc/services.d/crond/run
#!/usr/bin/with-contenv bash
crond -l 5 -f
